# supervisor/tracker_project_daphne.conf
[program:tracker_project_daphne]
user = vagrant
directory = /path/to/tracker_project/tracker_project
command = /path/to/tracker_project_venv/bin/daphne tracker_project.asgi:channel_layer
autostart = true
autorestart = true
stderr_logfile = /path/to/tracker_project/daphne_stderr.log
stdout_logfile = /path/to/tracker_project/daphne_stdout.log
stopsignal = INT
